# test-lab/workshop/profile.d/init.sh
#!/bin/bash

docker run -it -v .:/home/dev --cap-add=sys_admin sh3b0/recorder

