# labs/git/workshop/profile.d/init.sh
#!/bin/bash

export GIT_USERNAME=admin
export GIT_PASSWORD=admin
